# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true

[program:tailscale]
command=/bin/bash /workspaces/${LOCAL_WORKSPACE_FOLDER_BASENAME}/.devcontainer/tailscale.sh
autostart=true
autorestart=true
